# airootfs/etc/systemd/system/sddm.service.d/override.conf
[Service]
ExecStartPre=-/usr/bin/plymouth quit
